# docker-flask/supervisor.ini
[program:app]
command = /app/start-gunicorn.sh
directory=/app
user=alice
autostart=true
autorestart=true
